# service/fast_cpp_server.service
[Unit]
Description=Fast C++ Server Service
After=network.target

[Service]
# Run as root
User=root
Group=root

# Service type - simple because the start script waits on the process
Type=simple

# Prepare directories
ExecStartPre=/usr/bin/mkdir -p /var/fast_cpp_server/logs
ExecStartPre=/usr/bin/touch /var/fast_cpp_server/logs/fast_cpp_server.log
ExecStartPre=/usr/bin/chmod -R 755 /var/fast_cpp_server

# Start the wrapper script which will start and monitor the core program
ExecStart=/usr/local/bin/fast_cpp_server_dir/start-system.sh

# systemd handles restarts - always restart even on clean exit (exit code 0)
Restart=always
RestartSec=3s

# Working directory
WorkingDirectory=/usr/local/bin

# Runtime directory /run/fast_cpp_server
RuntimeDirectory=fast_cpp_server

# Allow write to /var (for logs)
ReadWritePaths=/var/fast_cpp_server
LogsDirectory=logs
LogsDirectoryMode=0755

# Resource limits
LimitNOFILE=65535
LimitCORE=infinity

# Set library path
Environment=LD_LIBRARY_PATH=/usr/local/lib/fast_cpp_server

[Install]
WantedBy=multi-user.target